# scripts/monitoring/fail2ban-metrics.timer
[Unit]
Description=Run Fail2ban metrics script every 1m

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m
AccuracySec=10s

[Install]
WantedBy=timers.target
