# scripts/monitoring/fail2ban-metrics.service
[Unit]
Description=Write Fail2ban Prometheus metrics

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fail2ban-metrics.sh
