# scripts/monitoring/ufw-metrics.service
[Unit]
Description=Write UFW Prometheus metrics

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ufw-metrics.sh
